# infra/server/base/dev/supervisord.conf
[program:tomcat]
command=/bin/bash /usr/local/apache-tomcat-8.0.22/bin/catalina.sh run

[program:h2]
command=java -jar h2-1.4.187.jar -webAllowOthers -tcpAllowOthers
directory=/usr/local/h2/bin
